# apps/udpping/others/udpping.service
[Unit]
Description=UDP Ping

[Service]
Type=simple
ExecStart=/opt/udpping

[Install]
WantedBy=multi-user.target